# roles/orchestrator_metrics/templates/orchestrator_metrics.timer.j2
[Unit]
Description=Run {{ project }} Orchestrator Metrics every 15 seconds

[Timer]
OnBootSec=5min
OnUnitActiveSec=15s
Persistent=true

[Install]
WantedBy=timers.target
